# ircddbgateway.timer
[Unit]
Description=ircDDBGateway Start delay timer

[Timer]
OnBootSec=30sec

[Install]
WantedBy=timers.target

