# ircddbgateway.timer
[Unit]
Description=ircDDBGateway Start delay timer (20210217)

[Timer]
OnBootSec=30sec

[Install]
WantedBy=timers.target

