# nextion.timer
[Unit]
Description = Nextion Service delay timer (20210217)

[Timer]
OnBootSec=40sec
#Unit=nextion.service

[Install]
WantedBy=timers.target
